# scripts/webapp.service
[Unit]
Description=CSYE6225 Web Application
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=csye6225
Group=csye6225
WorkingDirectory=/opt/csye6225
ExecStart=/usr/bin/node /opt/csye6225/server.js
Restart=always
RestartSec=10
StandardOutput=append:/opt/csye6225/logs/webapp.log
StandardError=append:/opt/csye6225/logs/webapp-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/csye6225

# Environment variables matching your MySQL setup
Environment="NODE_ENV=production"
Environment="PORT=8080"
Environment="DB_HOST=localhost"
Environment="DB_PORT=3306"
Environment="DB_NAME=health_check_db"
Environment="DB_USER=csye6225"
Environment="DB_PASS=MyPassword@123"

[Install]
WantedBy=multi-user.target