# scripts/webapp.service
[Unit]
Description=CSYE6225 Web Application
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=csye6225
Group=csye6225
WorkingDirectory=/opt/csye6225/webapp
EnvironmentFile=/opt/csye6225/webapp/.env
ExecStart=/usr/bin/node /opt/csye6225/webapp/server.js
Restart=always
RestartSec=10
StandardOutput=append:/opt/csye6225/logs/webapp.log
StandardError=append:/opt/csye6225/logs/webapp-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/csye6225

[Install]
WantedBy=multi-user.target